# raspcast.service
[Unit]
Description=Raspcast - Web-based mpv remote
After=network.target

[Service]
Type=simple
User=rupuin
WorkingDirectory=/home/rupuin/raspcast
EnvironmentFile=/etc/raspcast/env
Environment=HISTORY_STORE_PATH=/var/lib/raspcast/history.json
ExecStart=/home/rupuin/raspcast/raspcast
Restart=on-failure
KillMode=control-group
TimeoutStopSec=5
RestartSec=5

[Install]
WantedBy=multi-user.target
